# config/openrc/inittab
# /etc/inittab — jonerix system initialization table
#
# Format: id:runlevels:action:process
#
# Default runlevel
::sysinit:/bin/openrc sysinit
::sysinit:/bin/openrc boot
::wait:/bin/openrc default

# Spawn getty on tty1-tty3
tty1::respawn:/bin/getty 38400 tty1
tty2::respawn:/bin/getty 38400 tty2
tty3::respawn:/bin/getty 38400 tty3

# Serial console (for cloud / headless)
#ttyS0::respawn:/bin/getty -L 115200 ttyS0

# Shutdown / reboot
::shutdown:/bin/openrc shutdown
::ctrlaltdel:/bin/reboot
